# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ── TAP networking setup (runs first) ─────────────────────────────────────────
[program:tap-setup]
command=/usr/local/bin/tap-setup.sh
autostart=true
autorestart=false
startsecs=0
priority=5
stdout_logfile=/var/log/openqa/tap-setup.log
stderr_logfile=/var/log/openqa/tap-setup.log

# ── database ──────────────────────────────────────────────────────────────────
[program:postgresql]
command=/usr/bin/postgres -D /var/lib/pgsql/data
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql.log

# ── web server ────────────────────────────────────────────────────────────────
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/httpd-supervisor.log
stderr_logfile=/var/log/httpd-supervisor.log

# ── openQA services ───────────────────────────────────────────────────────────
[program:openqa-webui]
command=/usr/share/openqa/script/openqa daemon
user=geekotest
autostart=true
autorestart=true
priority=40
startsecs=5
stdout_logfile=/var/log/openqa/webui.log
stderr_logfile=/var/log/openqa/webui.log

[program:openqa-scheduler]
command=/usr/share/openqa/script/openqa-scheduler daemon
user=geekotest
autostart=true
autorestart=true
priority=50
startsecs=5
stdout_logfile=/var/log/openqa/scheduler.log
stderr_logfile=/var/log/openqa/scheduler.log

[program:openqa-websockets]
command=/usr/share/openqa/script/openqa-websockets daemon
user=geekotest
autostart=true
autorestart=true
priority=50
startsecs=5
stdout_logfile=/var/log/openqa/websockets.log
stderr_logfile=/var/log/openqa/websockets.log

[program:openqa-gru]
command=/usr/share/openqa/script/openqa gru run
user=geekotest
autostart=true
autorestart=true
priority=50
startsecs=5
stdout_logfile=/var/log/openqa/gru.log
stderr_logfile=/var/log/openqa/gru.log

# ── worker (single slot with TAP support) ─────────────────────────────────────
[program:openqa-worker]
command=/usr/local/bin/worker-wrapper.sh 1
user=geekotest
autostart=true
autorestart=true
priority=70
startsecs=15
startretries=3
stdout_logfile=/var/log/openqa/worker-1.log
stderr_logfile=/var/log/openqa/worker-1.log

# ── bootstrap (runs once to seed DB + API key) ────────────────────────────────
[program:bootstrap]
command=/usr/local/bin/bootstrap.sh
autostart=true
autorestart=false
startsecs=0
priority=65
stdout_logfile=/var/log/openqa/bootstrap.log
stderr_logfile=/var/log/openqa/bootstrap.log
